# contrib/systemd/bazzite-mcp-refresh.service
[Unit]
Description=Refresh bazzite-mcp docs cache
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/bazzite-mcp/.venv/bin/python -m bazzite_mcp.refresh
Environment=PATH=%h/.local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
